# packages/nix-single-user/system/nix-mkdir.service
[Unit]
DefaultDependencies=no
ConditionPathExists=!/nix
ConditionPathIsReadWrite=/
After=local-fs.target

[Service]
Type=oneshot
ExecStartPre=chattr -i /
ExecStart=mkdir /nix
ExecStopPost=chattr +i /
