# packages/nix-single-user/system/nix.mount
[Unit]
Description=Create /nix mount point.
DefaultDependencies=no
After=nix-mkdir.service
Wants=nix-mkdir.service

[Mount]
What=/var/nix
Where=/nix
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
